# opencode/opencode.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>opencode.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>fish</string>
    <string>-lc</string>
    <string>fnm exec --using=default opencode serve --hostname 0.0.0.0 --port 4096</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/bytedance</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>Crashed</key>
    <true/>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/opencode.server.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/opencode.server.err.log</string>
</dict>
</plist>
